# docker/train.cron
# Run the binary at midnight UTC
# 0 0 * * * /binaries/train.sh > /proc/1/fd/1 2>/proc/1/fd/2
* * * * * /binaries/train.sh /proc/1/fd/1 2>/proc/1/fd/2
